# systemd/hermit-crab-daemon.service
[Unit]
Description=Hermit Crab Daemon - 持续监控模式
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hermit_crab
EnvironmentFile=/root/hermit_crab/.env
ExecStart=/root/hermit_crab/venv/bin/python /root/hermit_crab/agent.py daemon

# 自动重启配置
Restart=always
RestartSec=60

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermit-crab-daemon

[Install]
WantedBy=multi-user.target

