# systemd/hermit-crab.service
[Unit]
Description=Hermit Crab - 寄居蟹自动迁移系统
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/hermit_crab
EnvironmentFile=/root/hermit_crab/.env
ExecStart=/root/hermit_crab/venv/bin/python /root/hermit_crab/agent.py check

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermit-crab

[Install]
WantedBy=multi-user.target

